# scripts/systemd/merge-media.service.example
[Unit]
Description=Mergerfs Mount Media
Wants=rclone-drive.service rclone-crypt.service
After=rclone-drive.service rclone-crypt.service
RequiresMountsFor=/media/drive /media/crypt /media/data

[Service]
Type=forking 
ExecStart=/usr/bin/mergerfs /media/data/local/media:/media/crypt/media /media/data/media \
-o rw,use_ino,allow_other,func.getattr=newest,category.action=all,category.create=ff,statfs_ignore=nc,cache.files=auto-full,func.getattr=newest,dropcacheonclose=true,nonempty
KillMode=process
Restart=on-failure
RestartSec=10
User=USER
Group=USER
#async_read=false,cache.files=partial,fsname=unionfs

[Install]
WantedBy=multi-user.target
